# Ethosight/install/templates/gpuserver.service.j2
[Unit]
Description=GPU Server
After=network.target

[Service]
ExecStart=/opt/scripts/run_gpuserver.sh
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target

