# Ethosight/install/templates/nginx.service.j2
[Unit]
Description=Nginx Service
After=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/bin/run_nginx.sh
Restart=always

[Install]
WantedBy=multi-user.target


